# WebProject/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=/usr/bin/java -jar /usr/share/jenkins/jenkins.war
autostart=true
autorestart=true
stderr_logfile=/var/log/jenkins.err.log
stdout_logfile=/var/log/jenkins.out.log

[program:spring]
command=java -jar /app/app.jar
autostart=true
autorestart=true
stderr_logfile=/var/log/spring.err.log
stdout_logfile=/var/log/spring.out.log